# distrib/build/sources/updates/update9/coretree.list
./etc/bashrc
./etc/dnsmasq.conf
./etc/httpd/conf/httpd.conf
./etc/rc.d/rc.updatered
./httpd/cgi-bin/advnet.cgi
./httpd/cgi-bin/apcupsd.cgi
./httpd/cgi-bin/backup.cgi
./httpd/cgi-bin/changepw.cgi
./httpd/cgi-bin/ddns.cgi
./httpd/cgi-bin/dhcp.cgi
./httpd/cgi-bin/dmzholes.cgi
./httpd/cgi-bin/extensions.cgi
./httpd/cgi-bin/hosts.cgi
./httpd/cgi-bin/ids.cgi
./httpd/cgi-bin/im.cgi
./httpd/cgi-bin/interfaces.cgi
./httpd/cgi-bin/ipblock.cgi
./httpd/cgi-bin/ipinfo.cgi
./httpd/cgi-bin/iptools.cgi
./httpd/cgi-bin/logs.cgi/firewalllog.dat
./httpd/cgi-bin/logs.cgi/ids.dat
./httpd/cgi-bin/logs.cgi/log.dat
./httpd/cgi-bin/logs.cgi/pop3log.dat
./httpd/cgi-bin/logs.cgi/proxylog.dat
./httpd/cgi-bin/logs.cgi/urlfilter.dat
./httpd/cgi-bin/modem.cgi
./httpd/cgi-bin/outgoing.cgi
./httpd/cgi-bin/p3scan.cgi
./httpd/cgi-bin/portfw.cgi
./httpd/cgi-bin/pppsetup.cgi
./httpd/cgi-bin/preferences.cgi
./httpd/cgi-bin/proxy.cgi
./httpd/cgi-bin/remote.cgi
./httpd/cgi-bin/shell.cgi
./httpd/cgi-bin/shutdown.cgi
./httpd/cgi-bin/sipproxy.cgi
./httpd/cgi-bin/smoothinfo.cgi
./httpd/cgi-bin/status.cgi
./httpd/cgi-bin/time.cgi
./httpd/cgi-bin/timedaccess.cgi
./httpd/cgi-bin/traffic.cgi
./httpd/cgi-bin/updates.cgi
./httpd/cgi-bin/urlfilter.cgi
./httpd/cgi-bin/vpnconn.cgi
./httpd/cgi-bin/vpnmain.cgi
./httpd/cgi-bin/xtaccess.cgi
./httpd/html/ui/css/style.css
./httpd/html/ui/img/banner.png
./httpd/html/ui/img/check-icon.png
./httpd/html/ui/img/frame-concrete/wbl.png
./httpd/html/ui/img/frame-concrete/wbr.png
./httpd/html/ui/img/frame-concrete/wbs.png
./httpd/html/ui/img/frame-concrete/wls.png
./httpd/html/ui/img/frame-concrete/wrs.png
./httpd/html/ui/img/frame-concrete/wtl.png
./httpd/html/ui/img/frame-concrete/wtr.png
./httpd/html/ui/img/frame-concrete/wts.png
./httpd/html/ui/img/x-icon.png
./usr/bin/smoothwall/rrdtool.pl
./usr/bin/smoothwall/trafficloader.pl
./usr/lib/smoothwall/header.pm
./usr/lib/smoothwall/langs/en.pl
./usr/lib/smoothwall/smoothtype.pm
./usr/lib/smoothwall/update.pm
./var/smoothwall/patches/installed
