# distrib/build/sources/updates/update8/coretree.list
./etc/group
./etc/bashrc
./etc/cron.daily/migrate-snort-pcap
./etc/ppp/adsl-cleanup
./etc/ppp/isdn-cleanup
./etc/rc.d/inc.block-ingress
./etc/rc.d/inc.links-up
./etc/rc.d/inc.open-ingress
./etc/rc.d/rc.firewall.up
./etc/rc.d/rc.netaddress.up
./etc/rc.d/rc.network
./etc/rc.d/rc.sysinit
./etc/rc.d/rc.updatered
./etc/rc.d/restorescript
./httpd/cgi-bin/backup.cgi
./httpd/cgi-bin/bandwidthbars.cgi
./httpd/cgi-bin/help.cgi
./httpd/cgi-bin/ids.cgi
./httpd/cgi-bin/interfaces.cgi
./httpd/cgi-bin/logs.cgi/firewalllog.dat
./httpd/cgi-bin/logs.cgi/ids.dat
./httpd/cgi-bin/logs.cgi/log.dat
./httpd/cgi-bin/logs.cgi/pop3log.dat
./httpd/cgi-bin/logs.cgi/proxylog.dat
./httpd/cgi-bin/outgoing.cgi
./httpd/cgi-bin/status.cgi
./httpd/cgi-bin/smoothinfo.cgi
./httpd/cgi-bin/time.cgi
./httpd/cgi-bin/updates.cgi
./httpd/html/ui/css/style.css
./httpd/html/ui/img/warning.png
./httpd/html/ui/js/script.js
./tmp/clamd
./usr/bin/installer/prepudevrules
./usr/bin/installer/installgrub
./usr/bin/smoothwall/addns.pl
./usr/bin/smoothwall/getx3systeminfo.pl
./usr/bin/smoothwall/imagerotator.pl
./usr/bin/smoothwall/machine_reg.pl
./usr/bin/smoothwall/ppp-off
./usr/bin/smoothwall/ppp-on
./usr/bin/smoothwall/rrdtool.pl
./usr/bin/smoothwall/smoothinfo.pl
./usr/bin/smoothwall/upsd-notify.pl
./usr/etc/sshd_config
./usr/lib/smoothwall/clamd.conf
./usr/lib/smoothwall/header.pm
./usr/lib/smoothwall/langs/alertboxes.en.pl
./usr/lib/smoothwall/langs/en.pl
./usr/lib/smoothwall/smoothtype.pm
./var/smoothwall/patches/installed
./var/smoothwall/proxy/cache
