# distrib/build/sources/updates/update2/coretree.list
./etc/cron.weekly/logrotate
./etc/cron.weekly/upddsttimes
./etc/httpd/conf/httpd.conf
./etc/logrotate.d/snort
./etc/logrotate.d/squid
./etc/rc.d/rc.sysinit
./etc/rc.d/rc.updatered
./etc/rc.d/restorescript
./httpd/cgi-bin/advnet.cgi
./httpd/cgi-bin/changepw.cgi
./httpd/cgi-bin/ddns.cgi
./httpd/cgi-bin/ipinfo.cgi
./httpd/cgi-bin/shutdown.cgi
./httpd/cgi-bin/smoothinfo.cgi
./httpd/cgi-bin/time.cgi
./httpd/cgi-bin/vpnconn.cgi
./usr/bin/installer/writegrubconf
./usr/bin/smoothwall/getdnsfromdhcpc.pl
./usr/bin/smoothwall/gettime.pl
./usr/bin/smoothwall/ppp-off
./usr/bin/smoothwall/smoothinfo.pl
./usr/bin/smoothwall/upddsttimes
./usr/bin/smoothwall/updkerneltz
./usr/bin/smoothwall/writentpd.pl
./usr/lib/smoothwall/header.pm
./usr/lib/smoothwall/langs/en.pl
./usr/lib/smoothwall/smoothtype.pm
./var/smoothwall/patches/installed
./var/smoothwall/proxy/acl
