# distrib/build/sources/updates/update3/coretree.list
./etc/apcupsd/apcupsd.conf
./etc/apcupsd/apcupsd.conf.template
./etc/inittab
./etc/logrotate.d/debug
./etc/logrotate.d/snort
./etc/miniupnpd.conf.template
./etc/passwd
./etc/rc.d/inc.block-ingress
./etc/rc.d/inc.links-up
./etc/rc.d/inc.open-ingress
./etc/rc.d/inc.rc-functions
./etc/rc.d/rc.adsl
./etc/rc.d/rc.firewall.up
./etc/rc.d/rc.halt
./etc/rc.d/rc.isdn
./etc/rc.d/rc.machineregister
./etc/rc.d/rc.netaddress.up
./etc/rc.d/rc.network
./etc/rc.d/rc.sysinit
./etc/rc.d/rc.updatered
./etc/syslog.conf
./httpd/cgi-bin/apcupsd.cgi
./httpd/cgi-bin/ddns.cgi
./httpd/cgi-bin/dhcp.cgi
./httpd/cgi-bin/dmzholes.cgi
./httpd/cgi-bin/graphs.cgi
./httpd/cgi-bin/help.cgi
./httpd/cgi-bin/hosts.cgi
./httpd/cgi-bin/index.cgi
./httpd/cgi-bin/ipblock.cgi
./httpd/cgi-bin/ipinfo.cgi
./httpd/cgi-bin/iptools.cgi
./httpd/cgi-bin/logs.cgi/proxylog.dat
./httpd/cgi-bin/logs.cgi/log.dat
./httpd/cgi-bin/outgoing.cgi
./httpd/cgi-bin/portfw.cgi
./httpd/cgi-bin/smoothinfo.cgi
./httpd/cgi-bin/time-clock.cgi
./httpd/cgi-bin/updates.cgi
./httpd/cgi-bin/urlfilter.cgi
./httpd/cgi-bin/xtaccess.cgi
./httpd/html/help/apcupsd.cgi.html.en
./httpd/html/ui/css/style.css
./httpd/html/ui/img/bear-teeth.png
./httpd/html/ui/img/urlfilter/adblock.gif
./httpd/html/ui/img/urlfilter/adblock.png
./httpd/html/ui/img/urlfilter/adblock.xcf
./httpd/html/ui/img/bottom-left.png
./httpd/html/ui/img/bottom-right.png
./httpd/html/ui/img/left.png
./httpd/html/ui/img/right.png
./httpd/html/ui/img/top-left.png
./httpd/html/ui/img/top-right.png
./httpd/html/ui/js/script.js
./usr/bin/smoothwall/apcupsdkillpower.pl
./usr/bin/smoothwall/upsd-notify.pl
./usr/bin/smoothwall/smoothinfo.pl
./usr/bin/smoothwall/stopsquid
./usr/bin/smoothwall/trafficloader.pl
./usr/bin/smoothwall/writeapcupsdconf.pl
./usr/bin/smoothwall/writeupnp
./usr/etc/dhcpcd.conf
./usr/lib/smoothwall/header.pm
./usr/lib/smoothwall/langs/alertboxes.en.pl
./usr/lib/smoothwall/langs/en.pl
./usr/lib/smoothwall/langs/glossary.en.pl
./usr/lib/smoothwall/menu/2000_Services/8000_ups.list
./usr/lib/smoothwall/smoothnet.pm
./usr/lib/smoothwall/smoothtype.pm
./var/smoothwall/apcupsd
./var/smoothwall/apcupsd/events
./var/smoothwall/apcupsd/settings
./var/smoothwall/proxy/cache
./var/smoothwall/proxy/squid.conf
