# distrib/build/sources/updates/update4/coretree.list
./etc/cron.weekly/logrotate
./etc/crontab
./etc/inittab
./etc/logrotate.base
./etc/logrotate.conf
./etc/logrotate.d/squid
./etc/rc.d/backupscript
./etc/rc.d/rc.sysinit
./etc/rc.d/restorescript
./httpd/cgi-bin/.htaccess
./httpd/cgi-bin/apcupsd.cgi
./httpd/cgi-bin/backup.cgi
./httpd/cgi-bin/ipinfo.cgi
./httpd/cgi-bin/register.cgi
./httpd/cgi-bin/status.cgi
./httpd/cgi-bin/traffic.cgi
./httpd/cgi-bin/trafficmonitor.cgi
./httpd/html/ui/js/script.js
./usr/bin/installer/prepudevrules
./usr/bin/smoothwall/gettime.pl
./usr/bin/smoothwall/writeproxy.pl
./usr/lib/smoothwall/langs/en.pl
./usr/lib/smoothwall/smoothtype.pm
./var/smoothwall/patches/installed
